# .config/systemd/user/name-resolv.service
[Unit]
Description = Name resolution test

[Service]
StandardError = null
StandardOutput = null

ExecStart = resolvectl query gnu.org
