# .config/systemd/user/name-resolv.target
[Unit]
Description = Name resolution

Requires = name-resolv.service

RefuseManualStart = 1
StopWhenUnneeded = 1
