# .config/systemd/user/update.timer
[Unit]
Description = Automatic package upgrades

[Timer]
OnBootSec = 45 minutes

[Install]
WantedBy = timers.target
